# darwin/Library/LaunchAgents/sh.opencode.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>sh.opencode.web</string>

	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/opencode</string>
		<string>web</string>
		<string>--port</string>
		<string>4096</string>
		<string>--hostname</string>
		<string>127.0.0.1</string>
		<string>--print-logs</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/alexramsay/git</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/alexramsay/.dprint/bin:/Users/alexramsay/.cargo/bin:/Users/alexramsay/go/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/alexramsay/Library/Logs/opencode-web.out.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/alexramsay/Library/Logs/opencode-web.err.log</string>
</dict>
</plist>
