# linux/home/.config/systemd/user/opencode.service
[Unit]
Description=OpenCode Web Server
After=network.target

[Service]
Type=simple
ExecStart=/home/linuxbrew/.linuxbrew/bin/opencode web
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
