# linux/.config/systemd/user/syncthing.service
[Unit]
Description=Syncthing

[Service]
ExecStart=/home/linuxbrew/.linuxbrew/bin/syncthing

[Install]
WantedBy=default.target
