# darwin/home/Library/LaunchAgents/sh.opencode.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>sh.opencode.web</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-c</string>
		<string>source ~/.zshrc &amp;&amp; exec /opt/homebrew/bin/opencode web --port 4096 --hostname 127.0.0.1 --print-logs</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/alexramsay/git</string>

	<key>RunAtLoad</key>
	<true/>

	<key>KeepAlive</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/alexramsay/Library/Logs/opencode-web.out.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/alexramsay/Library/Logs/opencode-web.err.log</string>
</dict>
</plist>
